# etc/systemd/system/domoticz.service.d/override.conf
[Service]
EnvironmentFile=-/etc/default/domoticz
ExecStart=
ExecStart=/opt/domoticz/domoticz -daemon -www $DOMOTICZ_WWW -sslwww $DOMOTICZ_SSLWWW